# crontab
0 13 * * * /go/tweet 0
0 7 * * * /go/tweet trend
0 23 * * * /go/tweet trend
